# source/syslog-backup/syslog-backup.cron
# Generated syslog backup schedule:
15 * * * * /usr/local/bin/syslog-backup.sh & > /dev/null

